# backend/crontab
# Send expiry warning emails daily at 9:00 AM UTC
0 9 * * * cd /app && python send_expiry_warnings.py
